# crontab
0 6 * * *  /app/run.sh
0 18 * * * /app/run.sh